# debian/prerm
#!/bin/sh

set -e

THAT_NAME=capsrvd

if [ "$1" = remove ]; then

systemctl is-active --quiet $THAT_NAME && \
    systemctl disable $THAT_NAME && \
    systemctl stop $THAT_NAME

fi


